# systemd/minecraft.socket
[Socket]
ListenFIFO=%t/minecraft.stdin
Service=minecraft.service