# systemd/minecraft.service
[Unit]
Description=Minecraft Server
Wants=network-online.target

[Service]
User=minecraft
WorkingDirectory=/opt/minecraft
ExecStart=/opt/minecraft/scripts/start_server.sh
ExecStop=/opt/minecraft/scripts/stop_server.sh
Sockets=minecraft.socket
StandardInput=socket

[Install]
WantedBy=multi-user.target
